# packaging/debian/postinst
#!/bin/sh
set -e

# Reload each active user's systemd manager and restart BlueCross services that
# are currently running, so an upgrade picks up the new binary automatically.
# try-restart is a no-op when a service is not running (e.g. fresh installs).
restart_running_user_services() {
    [ -d /run/user ] || return 0
    for rundir in /run/user/*; do
        [ -d "$rundir" ] || continue
        uid=$(basename "$rundir")
        user=$(id -nu "$uid" 2>/dev/null) || continue
        runuser -u "$user" -- env XDG_RUNTIME_DIR="$rundir" \
            systemctl --user daemon-reload >/dev/null 2>&1 || true
        for svc in bluecross-server.service bluecross-client.service; do
            runuser -u "$user" -- env XDG_RUNTIME_DIR="$rundir" \
                systemctl --user try-restart "$svc" >/dev/null 2>&1 || true
        done
    done
}

case "$1" in
    configure)
        # Reload udev rules to activate the new uinput rule
        if command -v udevadm > /dev/null 2>&1; then
            udevadm control --reload-rules || true
            udevadm trigger --subsystem-match=misc --attr-match=name=uinput || true
        fi

        # On upgrade ($2 is the previously-installed version), auto-restart
        # any running services with the new binary/units.
        if [ -n "$2" ]; then
            restart_running_user_services
        fi

        echo ""
        echo "BlueCross has been installed successfully!"
        echo ""
        echo "Post-installation steps:"
        echo "  1. Add your user to the 'input' group:"
        echo "     sudo usermod -aG input \$USER"
        echo ""
        echo "  2. Log out and back in (or reboot) for group changes to take effect"
        echo ""
        echo "  3. Copy the example config and set a shared 'psk' (>= 16 chars):"
        echo "     mkdir -p ~/.config/bluecross"
        echo "     cp /usr/share/bluecross/bluecross.json.example ~/.config/bluecross/bluecross.json"
        echo "     # generate a key:  head -c 32 /dev/urandom | base64"
        echo ""
        echo "  4. Edit the configuration file for your setup"
        echo ""
        echo "  5. (Optional) Enable the systemd user service:"
        echo "     systemctl --user daemon-reload"
        echo "     systemctl --user enable --now bluecross-server   # or bluecross-client"
        echo ""
        ;;
esac

#DEBHELPER#

exit 0
