# systemd/bluecross-server.service
[Unit]
Description=BlueCross Server - Keyboard/Mouse/Clipboard Sharing
Documentation=https://github.com/porcupin26/bluecross
After=network.target graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/bluecross server -f -c %h/.config/bluecross/bluecross.json
Restart=always
RestartSec=3

# Wayland-only. The user session manager normally exports WAYLAND_DISPLAY and
# XDG_RUNTIME_DIR into user services; if clipboard sharing cannot reach the
# compositor, import them explicitly, e.g.:
#   systemctl --user import-environment WAYLAND_DISPLAY XDG_RUNTIME_DIR
Environment=PATH=/usr/local/bin:/usr/bin

[Install]
WantedBy=graphical-session.target
