# systemd/bluecross-client.service
[Unit]
Description=BlueCross Client - Keyboard/Mouse Sharing
Documentation=https://github.com/yourusername/bluecross
After=network.target graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
WorkingDirectory=%h/projects/bluecross
ExecStart=%h/.local/bin/uv run bluecross-client -f -c %h/.config/bluecross/bluecross.json
Restart=on-failure
RestartSec=5

# Environment variables for display access
Environment=DISPLAY=:0
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin

[Install]
WantedBy=graphical-session.target
